# test/kea-supervisord.conf
[supervisord]
nodaemon=true

[program:kea-dhcp4]
command = kea-dhcp4 -c /usr/local/etc/kea/kea-dhcp4.conf
priority = 3
startinorder = true

[program:kea-dhcp6]
command = kea-dhcp6 -c /usr/local/etc/kea/kea-dhcp6.conf
priority = 4
startinorder = true

[program:kea-ctrl-agent]
command = kea-ctrl-agent -c /usr/local/etc/kea/kea-ctrl-agent.conf
priority = 5
startinorder = true

